# shell/profile.d/02-umask
#!/bin/sh
# 02-umask
# shellcheck enable=all
umask 7
