# shell/profile.d/12-:q
#!/bin/sh
# 12-:q
# shellcheck enable=all
alias :q='exit'
