# shell/profile.d/31-ps2
#!/bin/sh
# 31-ps2
# shellcheck enable=all
PS2="\${EMOJI}"
export PS2
