# config/shell/profile.d/01-misc
#!/bin/sh
# vim:ft=bash
# 01-common
# shellcheck enable=all
# shellcheck disable=SC3044
# shellcheck disable=SC3044
shopt -u cmdhist
shopt -s checkwinsize
umask 7
PATH="/npbin:${PATH}:${HOME}/.local/bin:${HOME}/.local/bin/*:."
export PATH
#shellcheck disable=SC2154
eval "LS_COLORS=\"st=30;4\${COL$(( 0 % COLORS ))}:no=3\${COL$(( 1 % COLORS ))}:fi=1;3\${COL$(( 1 % COLORS ))}:rs=0:di=3\${COL$(( 0 % COLORS ))};1:ln=3\${COL$(( 0 % COLORS ))};2:ex=1;3\${COL$(( 2 % COLORS ))}:\""
export LS_COLORS
WWW_HOME='http://lite.duckduckgo.com/lite/'
export WWW_HOME
