# config/shell/profile.d/03-path
#!/bin/sh
PATH="/npbin:${PATH}:${HOME}/.local/bin:${HOME}/.local/bin/*:."
export PATH
