# config/shell/profile.d/03-path
#!/bin/sh
# 03-path
# shellcheck enable=all
PATH="/npbin:${PATH}:${HOME}/.local/bin:${HOME}/.local/bin/*:."
export PATH
