# config/shell/profile.d/05-www-home
#!/bin/sh
WWW_HOME='http://lite.duckduckgo.com/lite/'
export WWW_HOME
