# config/shell/profile.d/11-btrls
#!/bin/sh
alias ls=btrls
btrls(){
	temp=""
	for x in "$@"; do
		t="$(echo "${x}" | sed -E 's|^-(\w*)a|-\1A|g;s| |\ |g;s|^-(\w*)l|-\1lh|g' || true)"
		[ -z "${temp}" ] && temp="${t}" || temp="${temp}ə${t}"
	done
	IFS="ə";
	#shellcheck disable=SC2086
	lsd ${_ls_extra:-} --color=auto ${temp}
}
export btrls
