# config/shell/profile.d/13-clear
#!/bin/sh
# 13-clear
# shellcheck enable=all
alias clear='printf "[H[3J";\clear;echo "[38;2;245;169;184mby using this terminal you agree to the TOS and Privacy Policy[0m"'
